# config/supervisor/devel/jupyterhub.ini
[program:jupyterhub]
command=/usr/local/bin/pixi run jupyterhub
stdout_logfile=/var/log/supervisor/jupyterhub.log
user=1000
directory=/workspace/pixi/devel/jupyter
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true

