# config/supervisor/devel/mcp-server-fetch.ini
[program:mcp-server-fetch]
command=/usr/local/bin/pixi run mcp-server-fetch
stdout_logfile=/var/log/supervisor/mcp-server-fetch.log
user=1000
directory=/workspace/pixi/devel/jupyter
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true
