# config/supervisor/devel/mcpo.ini
[program:mcpo]
command=/usr/local/bin/pixi run mcpo
stdout_logfile=/var/log/supervisor/mcpo.log
user=1000
directory=/workspace/pixi/devel/jupyter
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true
