# config/supervisor/devel/openwebui.ini
[program:openwebui]
command=pixi run openwebui
directory=/workspace/pixi/devel/openwebui
stdout_logfile=/var/log/supervisor/openwebui.log
user=1000
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true
