# config/supervisor/prod/mcp-server-time.ini
[program:mcp-server-time]
command=/usr/local/bin/pixi run mcp-server-time
stdout_logfile=/var/log/supervisor/mcp-server-time.log
user=1000
directory=/pixi/jupyter
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true
