# config/supervisor/prod/ollama.ini
[program:ollama]
command=/usr/bin/ollama serve
stdout_logfile=/var/log/supervisor/ollama.log
directory=/workspace
autostart=false
autorestart=true
user=1000
