# config/supervisor/prod/openwebui.ini
[program:openwebui]
command=/usr/local/bin/pixi run openwebui
directory=/pixi/openwebui
stdout_logfile=/var/log/supervisor/openwebui.log
user=1000
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true
