# config/jupyterhub.ini
[program:jupyterhub]
command=pixi run jupyterhub
user=1000
directory=/workspace/pixi/stable/jupyter
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true

