# config/jupyterhub-dev.ini
[program:jupyterhub-dev]
command=/usr/local/bin/pixi run jupyterhub-dev
user=1000
directory=/workspace/dev/pixi/jupyter
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true

