# config/jupyterbook.ini
[program:jupyterbook]
command=/usr/local/bin/pixi run jupyterbook
user=1000
directory=/workspace
autostart=false
autorestart=true
startsecs=2
priority=2

