# config/jupyterserver.ini
[program:jupyterserver]
command=/usr/local/bin/pixi run jupyterserver
user=1000
directory=/workspace/pixi/jupyter
autostart=false
autorestart=true
startsecs=2
priority=2
stopasgroup=true
killasgroup=true

