# config/jupyterhub.ini
[program:jupyterhub]
command=/usr/local/bin/pixi run jupyterhub
user=root
directory=/workspace
autostart=false
autorestart=true
startsecs=2
priority=2

