# config/jupyterbook.ini
[program:jupyterbook]
command=/usr/local/bin/pixi run jupyterbook
user=atrawog
directory=/workspace
autostart=false
autorestart=true
startsecs=2
priority=2

