# config/ollama.ini
[program:ollama]
command=/usr/bin/ollama serve
autostart=false
autorestart=true
user=1000
