# config/openwebui-uv.ini
[program:openwebui-uv]
command=/usr/sbin/uv run open-webui serve
directory=/workspace
autostart=false
autorestart=true
user=1000
stopasgroup=true
killasgroup=true
