# config/openwebui-dev.ini
[program:openwebui]
command=/usr/local/bin/pixi run openwebui-dev
user=1000
directory=/workspace/pixi/openwebui-dev
autostart=false
autorestart=true
startsecs=2
priority=2